# glediator/supervisor/conf.d/glediator.conf
[program:glediator]
environment=DISPLAY=:1
command=java -jar /usr/src/app/Glediator_V2.jar