# glediator/supervisor/conf.d/x11vnc.conf
[program:x11vnc]
command=x11vnc -display :1 -forever -nopw
