# glediator/supervisor/conf.d/xvfb.conf
[program:xvfb]
command=Xvfb :1 -screen 0 1280x720x24
